# airootfs/etc/systemd/system/intel.service
[Unit]
Description=Detect legacy Intel GPU (Gen 1-4) and install xf86-video-intel
DefaultDependencies=no
After=local-fs.target
Before=multi-user.target graphical.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/intel-legacy-gpu-installer
StandardOutput=journal
StandardError=journal
MemoryMax=256M

[Install]
WantedBy=multi-user.target
