# airootfs/etc/systemd/system/getty@tty1.service.d/autologin.conf
[Service]
ExecStartPre=-/bin/sleep 1
ExecStart=
ExecStart=-/sbin/agetty -o '-p -f -- \\u' --noclear --autologin liveuser - $TERM
